# data/etc/inittab
::respawn:/bin/sh /etc/init.d/accelerometer
::respawn:/bin/sh /etc/init.d/nmea
